# linux/joke.timer
[Unit]
Description=Ask for a joke every two hours

[Timer]
OnCalendar=0/2:00:00
Unit=joke.service

[Install]
WantedBy=timers.target
