# linux/joke.service
[Unit]
Description=Tell me a joke

[Service]
Type=oneshot
ExecStart=/home/user/.config/systemd/user/joke.sh
